# ollama-kie-proxy.service
[Unit]
Description=Ollama-KIE.AI Proxy Service
Documentation=https://github.com/yourusername/ollama-kie-proxy
After=network.target

[Service]
Type=simple
User=your_username
WorkingDirectory=/path/to/copilot-ollama-kie-proxy
EnvironmentFile=/path/to/copilot-ollama-kie-proxy/.env
ExecStart=/usr/bin/python3 -m uvicorn main:app --host 127.0.0.1 --port 11434 --workers 4
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
